# flexget/root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc flexget daemon start -d -c /flexget/config.yml